# systemd/.config/systemd/user/batterybrightness.timer
[Unit]
Description=Checks batterybrightness every 10mins

[Timer]
OnCalendar=*-*-* *:0..59/10:00

[Install]
WantedBy=timers.target

