# systemd/.config/systemd/user/timers.target.wants/change_wallpaper.timer
[Unit]
Description=Changes wallpaper every hour

[Timer]
OnCalendar=*-*-* *:00:*

[Install]
WantedBy=timers.target
